# service/datastax-agent/run
#!/usr/bin/env bash

echo "stomp_interface: $OPS_IP" | sudo tee -a /var/lib/datastax-agent/conf/address.yaml

/usr/share/datastax-agent/bin/datastax_agent_monitor >> /var/log/datastax-agent/agent.log 2>&1

exec 2>&1
exec setuidgid opscenter-agent --exec /usr/share/datastax-agent/bin/datastax-agent
